# docs/systemd/ai-usage-collector.service
[Unit]
Description=AI usage status collector
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/ai-usage-collect
