# systemd/ai-usage-collector.timer
[Unit]
Description=Run AI usage collector every 2 minutes

[Timer]
OnBootSec=30
OnUnitActiveSec=120
Unit=ai-usage-collector.service

[Install]
WantedBy=timers.target
