# systemd/ai-usage-collector.timer
[Unit]
Description=Run AI usage collector every 3 minutes

[Timer]
OnBootSec=30
OnUnitActiveSec=180
Unit=ai-usage-collector.service

[Install]
WantedBy=timers.target
